# extras/com.job-agent.cloudflared-named.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Named Cloudflare tunnel (stable URL). Set CLOUDFLARE_TUNNEL_TOKEN in .env first.
  Install: ./scripts/install-mac-remove-stack.sh
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.job-agent.cloudflared-named</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/arkadiykats/devops-job-agent</string>
  <key>StandardOutPath</key>
  <string>/Users/arkadiykats/logs/digest-remove-cloudflared.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/arkadiykats/logs/digest-remove-cloudflared.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/arkadiykats</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/cloudflared</string>
    <string>tunnel</string>
    <string>run</string>
    <string>--token</string>
    <string>CLOUDFLARE_TUNNEL_TOKEN_PLACEHOLDER</string>
  </array>
</dict>
</plist>
